# updatedb/systemd/updatedb.service
[Unit]
Description=Update the mlocate database

[Service]
Type=oneshot
ExecStart=/usr/bin/updatedb
